# modules/crawler/crontab
00 00 * * * cd /src && /usr/local/bin/python crone_update_script.py
# Mandatory blank line